# deploy/house-dashboard.service
[Unit]
Description=House dashboard server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/house
EnvironmentFile=/opt/house/.env
ExecStart=/opt/house/.venv/bin/house dashboard --host 127.0.0.1 --port 8765
Restart=always
RestartSec=10
User=house
Group=house

[Install]
WantedBy=multi-user.target
